# server_image/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:fcatalog_server]
command=/bin/bash -c "sudo -u ufcatalog /home/ufcatalog/ufcatalog_env/bin/python /home/ufcatalog/bin/fcatalog_server 0.0.0.0 1337"
